# docker/workspace/crontab/laradock
* * * * * laradock /usr/bin/php yii schedule/run --scheduleFile=@console/cron/schedule.php 1>> /dev/null 2>&1

